# crons/crontab
# m h  dom mon dow user command
  * *   *  *   *   default  ls
